# scripts/launchd/com.gurusishya.drive-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Sync staging folder to Google Drive every 6 hours -->
    <key>Label</key>
    <string>com.gurusishya.drive-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npx</string>
        <string>tsx</string>
        <string>scripts/sync-to-drive.ts</string>
        <string>--sync</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/racit/PersonalProject/video-pipeline</string>

    <!-- Every 6 hours: 6AM, 12PM, 6PM, midnight -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>18</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/racit/guru-sishya-logs/drive-sync-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/racit/guru-sishya-logs/drive-sync-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/racit</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>15</integer>
</dict>
</plist>
