# scripts/launchd/com.gurusishya.auto-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
      Auto-publish: Tue/Thu/Sat at 7:15 PM IST (= 1:45 PM UTC = local depends on timezone).
      macOS uses local time for StartCalendarInterval.
      IST = UTC+5:30, so 19:15 IST = 13:45 UTC.
      Adjust Hour/Minute below to your macOS local timezone.
      For IST (if macOS is set to IST): Hour=19, Minute=15
      For UTC: Hour=13, Minute=45
    -->
    <key>Label</key>
    <string>com.gurusishya.auto-publish</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npx</string>
        <string>tsx</string>
        <string>scripts/auto-publish.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/racit/PersonalProject/video-pipeline</string>

    <!-- Tuesday=2, Thursday=4, Saturday=6 at 19:15 IST -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>19</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>19</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>19</integer>
            <key>Minute</key>
            <integer>15</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/racit/guru-sishya-logs/auto-publish-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/racit/guru-sishya-logs/auto-publish-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/racit</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Restart on failure (non-zero exit) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>600</integer>
</dict>
</plist>
